# services/testService2/run
#!/bin/bash

echo testService2
echo hi
sleep 8
echo bye
sleep 1
exit 0
