# services/testService2/run
#!/bin/bash

echo testService2 PID: $$
sleep infinity